# dist/binaries/services/dev.querybird.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>dev.querybird</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/querybird</string>
        <string>start</string>
        <string>--config-dir</string>
        <string>/opt/querybird/configs</string>
        <string>--log-level</string>
        <string>info</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/opt/querybird/logs/querybird.log</string>
    <key>StandardErrorPath</key>
    <string>/opt/querybird/logs/querybird.error.log</string>
    <key>WorkingDirectory</key>
    <string>/opt/querybird</string>
</dict>
</plist>
